# rc.d/fed
#!/bin/ksh
daemon="/usr/local/bin/fed"
daemon_flags=""
rc_bg=YES
. /etc/rc.d/rc.subr
rc_cmd "$1"
