# rc.d/listener
#!/bin/ksh
daemon="/usr/local/bin/listener"
daemon_flags=""
rc_bg=YES
. /etc/rc.d/rc.subr
rc_cmd "$1"
